# systemd/zatsuma-duckdns.timer
[Unit]
Description=Zatsuma Timer (DuckDNS)

[Timer]
OnUnitActiveSec=2min
Unit=zatsuma-duckdns.service

