# systemd/zatsuma-duckdns.service
[Unit]
Description="Zatsuma DuckDNS (Dynamic IP)"

[Service]
Type=oneshot
ExecStart=/usr/local/bin/zatsuma/shopd -duckdns
User=USERNAME

