# systemd/zatsuma-prices.timer
[Unit]
Description=Zatsuma Timer (Exchange Rates)

[Timer]
OnUnitActiveSec=1min
Unit=zatsuma-prices.service

