# systemd/zatsuma-prices.service
[Unit]
Description="Zatsuma Exchange Rate"

[Service]
Type=oneshot
ExecStart=/usr/local/bin/zatsuma/shopd -coinlib
User=USERNAME

